# templates/default/systemd/imply.service.erb
[Unit]
Description=Imply - <%= @type %> <%= @service %>
After=network.target

[Service]
Type=simple
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @prefix_root %>/imply/bin/run-<%= @type %> <%= @service %>
Restart=on-failure

[Install]
WantedBy=multi-user.target
