# templates/default/systemd/pivot.service.erb
[Unit]
Description=Imply Pivot
After=network.target

[Service]
Type=fork
User=<%= @user %>
Group=<%= @group %>
ExecStart=/usr/bin/pivot -p <%= @pivot_port %> --druid <%= @druid_broker %>

[Install]
WantedBy=multi-user.target
